# config/crontab
@reboot su server -c "/home/server/framework/script/init.sh"
0 0 * * * su server -c "/home/server/framework/script/reset.sh" && reboot
