# config/crontab
@reboot su server -c '~/framework/script/init.sh'
0 0 * * * su server -c '~/framework/script/reset.sh' && reboot
